# scripts/daemon/crontab
# minecraft
0 4 * * * /usr/bin/python3 /home/projects/mc/minecraft-zerbitzaria/scripts/backup/backup.py >> /dev/null 2>> /home/projects/mc/minecraft-zerbitzaria/scripts/backup/backup.log